# root/etc/s6-overlay/s6-rc.d/init-roombarr-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

mkdir -p /config/logs
lsiown -R abc:abc /config/logs
lsiown abc:abc /app
