# root/etc/s6-overlay/s6-rc.d/init-roombarr-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

mkdir -p /data
lsiown -R abc:abc /data
lsiown abc:abc /app
